# 02-ubuntu-nginx-java/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx/nginx.log
stderr_logfile=/var/log/nginx/nginx_error.log
priority=10

[program:http8081]
command=java -jar /usr/sbin/http-server.jar 8081
autostart=true
autorestart=true
priority=11

[program:http8082]
command=java -jar /usr/sbin/http-server.jar 8082
autostart=true
autorestart=true
priority=12